# src/etc/cron.d/dyndns
* * * * * root /cf-ddns.sh >/var/log/dns.log 2>&1


